# ansible/roles/mounts/templates/backup-vps-proxy.timer.j2
[Unit]
Description=Run vps-proxy backup nightly
Requires=backup-vps-proxy.service

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
